# kin-office.service
[Unit]
Description=Kin Office Module - Nextcloud with OnlyOffice
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Pulls / healthchecks can exceed the default start timeout; image first-pull can take many minutes.
TimeoutStartSec=2400
ExecStart=/opt/kin/modules/kin-office/kin-office-wrapper.sh
ExecStop=/opt/kin/modules/kin-office/kin-office-wrapper.sh stop
ExecReload=/opt/kin/modules/kin-office/kin-office-wrapper.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
